# profiles/desktop/airootfs/etc/systemd/system/intel.service
[Unit]
Description=Detect and configure Intel legacy GPU drivers
DefaultDependencies=no
After=local-fs.target systemd-modules-load.service
Before=display-manager.service graphical.target
Wants=systemd-modules-load.service
ConditionPathExists=/usr/bin/intel-legacy-gpu
ConditionVirtualization=!container

[Service]
Type=oneshot
ExecStart=/usr/bin/intel-legacy-gpu
RemainAfterExit=yes
TimeoutSec=30
StandardOutput=journal
StandardError=journal

# Security hardening
PrivateDevices=no
ProtectSystem=false
ProtectHome=true
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
