# archiso/airootfs/etc/skel/.config/autostart/welcome-center.desktop
[Desktop Entry]
Name=Welcome Center
Comment=Launches the CrystalLinux Welcome Center
Exec=/usr/bin/welcome-center
Icon=CrystalLinux-logo
Terminal=false
Type=Application
Categories=System;Utility;
StartupNotify=false
X-GNOME-Autostart-enabled=true
